# _config/systemd/user/kill-ssh-sessions.service
[Unit]
Description=Kill all ssh sessions before system turns off
Before=network.target

[Service]
Type=simple
ExecStart=-/usr/bin/sh %h/.config/systemd/user/scripts/kill-ssh-sessions

[Install]
WantedBy=halt.target poweroff.target reboot.target shutdown.target sleep.target
